# swell-arb.service
[Unit]
Description=Swell Fee Flow Arbitrage Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/swell-auction
EnvironmentFile=/home/pi/swell-auction/.env
ExecStart=/home/pi/swell-auction/venv/bin/python bot/main.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swell-arb

[Install]
WantedBy=multi-user.target
