# swell-arb.service
[Unit]
Description=Swell Fee Flow Arbitrage Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fenyvesi
WorkingDirectory=/home/fenyvesi/Projects/swell-auction
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/fenyvesi/Projects/swell-auction/venv/bin/python /home/fenyvesi/Projects/swell-auction/bot/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
